# install/gunicorn_webchat.socket
[Unit]
Description=gunicorn socket

[Socket]
ListenStream=/run/gunicorn/webchat.socket

[Install]
WantedBy=sockets.target